# deploy/systemd/pairly-bot.service
# systemd unit for the Pairly Telegram bot (aiogram polling).
# Install: /etc/systemd/system/pairly-bot.service  then `systemctl daemon-reload`.
#
# Runs as a dedicated unprivileged `pairly` user. Secrets come from
# /etc/pairly/pairly.env (chmod 600, owner pairly). See deploy/scripts/install.sh.
[Unit]
Description=Pairly Telegram bot (aiogram polling)
Documentation=https://github.com/PAIRLY/pairly/blob/main/docs/runbook.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pairly
Group=pairly

# Repo + venv. Created by install.sh.
WorkingDirectory=/opt/pairly
EnvironmentFile=/etc/pairly/pairly.env

# The package is installed into the venv by `uv sync`, so `pairly.*` imports resolve.
ExecStart=/opt/pairly/.venv/bin/python -m pairly.main

Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# Narrow write allowlist. Source code under /opt/pairly (repo checkout) is
# read-only in spirit: it is deployed from git and never mutated at runtime.
# The only runtime-mutable surfaces are the SQLite DB (/var/lib/pairly), the
# log dir (/var/log/pairly), and the venv's __pycache__ (/opt/pairly/.venv).
ReadWritePaths=/var/lib/pairly /var/log/pairly /opt/pairly/.venv
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true

# Drop all capabilities (process runs unprivileged).
CapabilityBoundingSet=
AmbientCapabilities=
# AF_UNIX is needed for journald / sd_notify sockets; AF_INET/AF_INET6 for HTTP.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
# Restrict to the curated @system-service syscall set; deny others with EPERM.
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectProc=invisible
UMask=0027

# Journald captures stdout/stderr; also tag for readability.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pairly-bot

[Install]
WantedBy=multi-user.target
